# backend/app/db/seeds/mitre_attack.json
[
  {"id": "T1595", "name": "Active Scanning", "tactic": "reconnaissance", "description": "Adversaries scan victim IP blocks or services to gather information for targeting."},
  {"id": "T1592", "name": "Gather Victim Host Information", "tactic": "reconnaissance", "description": "Adversaries gather details about victim hosts to plan follow-on activity."},
  {"id": "T1583", "name": "Acquire Infrastructure", "tactic": "resource-development", "description": "Adversaries buy, lease, or rent infrastructure for use throughout the attack lifecycle."},
  {"id": "T1588", "name": "Obtain Capabilities", "tactic": "resource-development", "description": "Adversaries buy or steal capabilities such as malware, exploits, or certificates."},
  {"id": "T1566", "name": "Phishing", "tactic": "initial-access", "description": "Adversaries send phishing messages to gain access to victim systems."},
  {"id": "T1190", "name": "Exploit Public-Facing Application", "tactic": "initial-access", "description": "Adversaries exploit a weakness in an internet-facing host or application."},
  {"id": "T1078", "name": "Valid Accounts", "tactic": "initial-access", "description": "Adversaries use compromised credentials to bypass access controls."},
  {"id": "T1133", "name": "External Remote Services", "tactic": "initial-access", "description": "Adversaries abuse external remote services (VPN, Citrix, RDP gateways) to gain initial access."},
  {"id": "T1059", "name": "Command and Scripting Interpreter", "tactic": "execution", "description": "Adversaries abuse command/script interpreters to execute commands or payloads."},
  {"id": "T1059.001", "name": "PowerShell", "tactic": "execution", "description": "Adversaries abuse PowerShell for execution, discovery, and lateral movement."},
  {"id": "T1203", "name": "Exploitation for Client Execution", "tactic": "execution", "description": "Adversaries exploit software vulnerabilities in client applications to execute code."},
  {"id": "T1204", "name": "User Execution", "tactic": "execution", "description": "Adversaries rely on a user to execute a malicious file or link."},
  {"id": "T1547", "name": "Boot or Logon Autostart Execution", "tactic": "persistence", "description": "Adversaries configure system settings to execute a program during system boot or logon."},
  {"id": "T1053", "name": "Scheduled Task/Job", "tactic": "persistence", "description": "Adversaries abuse task scheduling to execute code for persistence or execution."},
  {"id": "T1136", "name": "Create Account", "tactic": "persistence", "description": "Adversaries create accounts to maintain access to victim systems."},
  {"id": "T1098", "name": "Account Manipulation", "tactic": "persistence", "description": "Adversaries manipulate accounts to maintain access, e.g. adding credentials or permissions."},
  {"id": "T1068", "name": "Exploitation for Privilege Escalation", "tactic": "privilege-escalation", "description": "Adversaries exploit software vulnerabilities to elevate privileges."},
  {"id": "T1055", "name": "Process Injection", "tactic": "privilege-escalation", "description": "Adversaries inject code into processes to evade defenses and/or elevate privileges."},
  {"id": "T1548", "name": "Abuse Elevation Control Mechanism", "tactic": "privilege-escalation", "description": "Adversaries bypass mechanisms designed to control elevated privileges (e.g. UAC)."},
  {"id": "T1070", "name": "Indicator Removal", "tactic": "defense-evasion", "description": "Adversaries delete or alter artifacts to remove evidence of their presence."},
  {"id": "T1027", "name": "Obfuscated Files or Information", "tactic": "defense-evasion", "description": "Adversaries obfuscate content to make it difficult to discover or analyze."},
  {"id": "T1036", "name": "Masquerading", "tactic": "defense-evasion", "description": "Adversaries manipulate features of artifacts to make them appear legitimate."},
  {"id": "T1562", "name": "Impair Defenses", "tactic": "defense-evasion", "description": "Adversaries disable or modify security tools to avoid detection."},
  {"id": "T1110", "name": "Brute Force", "tactic": "credential-access", "description": "Adversaries use brute force techniques to gain access to accounts."},
  {"id": "T1110.001", "name": "Password Guessing", "tactic": "credential-access", "description": "Adversaries guess passwords without prior knowledge of leaked/plaintext credentials."},
  {"id": "T1003", "name": "OS Credential Dumping", "tactic": "credential-access", "description": "Adversaries dump credentials from the OS to obtain account login information."},
  {"id": "T1552", "name": "Unsecured Credentials", "tactic": "credential-access", "description": "Adversaries search for insecurely stored credentials such as config files or scripts."},
  {"id": "T1556", "name": "Modify Authentication Process", "tactic": "credential-access", "description": "Adversaries modify authentication mechanisms to access credentials or bypass MFA."},
  {"id": "T1082", "name": "System Information Discovery", "tactic": "discovery", "description": "Adversaries gather detailed information about the OS and hardware configuration."},
  {"id": "T1087", "name": "Account Discovery", "tactic": "discovery", "description": "Adversaries enumerate local or domain accounts on a system."},
  {"id": "T1018", "name": "Remote System Discovery", "tactic": "discovery", "description": "Adversaries enumerate other systems on a network for lateral movement targeting."},
  {"id": "T1046", "name": "Network Service Discovery", "tactic": "discovery", "description": "Adversaries scan for network services running on remote hosts."},
  {"id": "T1021", "name": "Remote Services", "tactic": "lateral-movement", "description": "Adversaries use valid accounts to log into services designed for remote access."},
  {"id": "T1021.001", "name": "Remote Desktop Protocol", "tactic": "lateral-movement", "description": "Adversaries use RDP to move laterally within a network."},
  {"id": "T1550", "name": "Use Alternate Authentication Material", "tactic": "lateral-movement", "description": "Adversaries use stolen tokens, hashes, or tickets to bypass normal auth."},
  {"id": "T1005", "name": "Data from Local System", "tactic": "collection", "description": "Adversaries search local system sources to find files of interest."},
  {"id": "T1114", "name": "Email Collection", "tactic": "collection", "description": "Adversaries target user email to collect sensitive information."},
  {"id": "T1560", "name": "Archive Collected Data", "tactic": "collection", "description": "Adversaries compress and/or encrypt collected data prior to exfiltration."},
  {"id": "T1071", "name": "Application Layer Protocol", "tactic": "command-and-control", "description": "Adversaries use application-layer protocols to blend C2 traffic with legitimate traffic."},
  {"id": "T1105", "name": "Ingress Tool Transfer", "tactic": "command-and-control", "description": "Adversaries transfer tools or files from an external system into a compromised environment."},
  {"id": "T1573", "name": "Encrypted Channel", "tactic": "command-and-control", "description": "Adversaries employ encryption to conceal C2 traffic from network defenses."},
  {"id": "T1090", "name": "Proxy", "tactic": "command-and-control", "description": "Adversaries use a connection proxy to direct C2 traffic and mask its true source."},
  {"id": "T1041", "name": "Exfiltration Over C2 Channel", "tactic": "exfiltration", "description": "Adversaries steal data by exfiltrating it over an existing C2 channel."},
  {"id": "T1048", "name": "Exfiltration Over Alternative Protocol", "tactic": "exfiltration", "description": "Adversaries exfiltrate data using a protocol different from the existing C2 channel."},
  {"id": "T1486", "name": "Data Encrypted for Impact", "tactic": "impact", "description": "Adversaries encrypt data to interrupt availability, e.g. ransomware."},
  {"id": "T1490", "name": "Inhibit System Recovery", "tactic": "impact", "description": "Adversaries delete or disable OS recovery features to prevent restoration."},
  {"id": "T1489", "name": "Service Stop", "tactic": "impact", "description": "Adversaries stop or disable services to render systems unavailable."},
  {"id": "T1498", "name": "Network Denial of Service", "tactic": "impact", "description": "Adversaries degrade or block availability of targeted resources to users."}
]
